# com.wulala.claude-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.wulala.claude-bridge</string>
  <key>Program</key><string>/bin/bash</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>exec /Users/wulala/Downloads/Helper/claude-bridge/bridge/bridge_launch.sh</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/wulala/Downloads/Helper/claude-bridge/bridge</string>
  <key>KeepAlive</key><true/>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/tmp/com.wulala.claude-bridge.stdout.log</string>
  <key>StandardErrorPath</key><string>/tmp/com.wulala.claude-bridge.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/Users/wulala/.npm-global/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key><string>/Users/wulala</string>
    <key>BRIDGE_PORT</key><string>8766</string>
    <key>BRIDGE_DISABLE_MDNS</key><string>1</string>
  </dict>
</dict>
</plist>
